# docker/services/app/resources/etc/supervisor/conf.d/varnish.conf
[program:varnish]
command=varnishd -f /etc/varnish/default.vcl -a 0.0.0.0:80 -s malloc,1G -F
priority=998
user=root
autostart=true
autorestart=true
startretries=3
